# cabinserver/conf/supervisor.conf
[program:cabinserver]
user=root
command=/srv/cabinserver
autostart=true
autorestart=true
stderr_logfile=/var/log/cabinserver.err.log
stdout_logfile=/var/log/cabinserver.out.log

